# s6/cron/run
#!/bin/sh

set -eux

exec crond -f -L /dev/stdout
